# sys/kaplanbot.service
[Unit]
Description=KaplanBots startup 

[Service]
ExecStart=/bin/bash /usr/bin/kaplanbot.sh

[Install]
WantedBy=multi-user.target

